# ruijie-gdstvc-autologin/systemd/ruijie.service
[Unit]
Description=Ruijie Network Auto-Login Service
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/var/run/ruijie-daemon.pid
ExecStart=/usr/local/bin/ruijie.sh --daemon
ExecStop=/usr/local/bin/ruijie.sh --stop
Restart=on-failure
RestartSec=10

# Security hardening
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
NoNewPrivileges=true

# Run as unprivileged user if possible
# User=ruijie
# Group=ruijie

[Install]
WantedBy=multi-user.target
